# units/backup-3.timer
[Unit]
Description=Run backup

[Timer]
OnCalendar=daily
Persistent=true
#WakeSystem=true

[Install]
WantedBy=default.target
