# units/backup-1.timer
[Unit]
Description=Run backup

[Timer]
OnBootSec=10min
OnUnitActiveSec=1w

[Install]
WantedBy=default.target
