# units/forkbomb-1.service
[Service]
MemoryLimit=30M
CPUQuota=20%
OOMScoreAdjust=1000
ExecStart=/usr/bin/perl -e "fork while 1"
