# units/ifork-1.service
[Service]
Type=forking
ExecStart=/usr/local/bin/ifork

[Install]
WantedBy=multi-user.target
