# units/socat-1.service
[Service]
ExecStart=/usr/bin/socat TCP-LISTEN:8888 'SYSTEM:echo Hello World'
