# units/socat-2.service
[Service]
ExecStart=/usr/bin/socat TCP-LISTEN:8888,reuseaddr 'SYSTEM:echo Hello World'
Restart=on-success
